# systemd/keycloak.service
# keycloak.service — the Keycloak server (ADR-0005).
# Pre-built at bake; starts with --optimized (ADR-0004). Runs unprivileged.
[Unit]
Description=Keycloak
After=network-online.target keycloak-config.service
Wants=network-online.target
Requires=keycloak-config.service

[Service]
Type=exec
User=keycloak
Group=keycloak
# The neutral keycloak.conf is baked into KEYCLOAK_HOME/conf and read natively
# (ADR-0002); no KC_CONFIG_FILE needed. Env + secrets are boot-injected onto
# tmpfs (ADR-0008) by keycloak-config.service (bootstrap-admin creds, when set,
# are routed into secrets.env); '-' = optional.
EnvironmentFile=-/run/keycloak/keycloak.env
EnvironmentFile=-/run/keycloak/secrets.env
ExecStart=/opt/keycloak/bin/kc.sh start --optimized
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

# Hardening (defense-in-depth, subordinate to SELinux — ADR-0005/0011)
# Keycloak runs unprivileged and owns its own KEYCLOAK_HOME/data, which it must
# write (gzip cache, transaction logs). We deliberately do NOT use
# ProtectSystem=strict: it made the install tree read-only and fought Keycloak's
# normal need to write under its home, which required a fragile symlink dance
# (keycloak#31949). Security here rests on the unprivileged user, SELinux
# Enforcing, the immutable/neutral AMI, and the single-purpose node.
NoNewPrivileges=yes
ProtectHome=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
